# deploy/exeunt-alert-responder.path
[Unit]
Description=Watch for new alert emails

[Path]
DirectoryNotEmpty=/home/exedev/Maildir/new
MakeDirectory=yes

[Install]
WantedBy=multi-user.target
